# znc/root/etc/service/znc/run
#!/bin/sh
exec s6-setuidgid znc znc --datadir /etc/znc --foreground
