# prosody/etc/service/prosody/run
#!/bin/sh
export __FLUSH_LOG=yes
exec s6-setuidgid prosody prosody
