# pritunl/etc/service/mongodb/run
#!/bin/sh
mkdir -p /data/db
exec mongod --dbpath /data/db
