# seafile/etc/service/seafile/finish
#!/bin/sh
if [ $1 -ne 0 -a $1 -ne 256 ]; then
    echo exec halt
fi
