# seafile/etc/service/seafile/run
#!/bin/sh
exec seafile-controller -f -c $CCNET_CONF_DIR -d $SEAFILE_CONF_DIR -F $SEAFILE_CENTRAL_CONF_DIR
