# piwik/root/etc/periodic/hourly/piwik-archive.sh
#!/bin/sh
exec php7 /srv/www/console core:archive
