# piwik/root/etc/periodic/hourly/piwik-archive.sh
#!/bin/sh
exec php7 /var/www/console core:archive
