# base/root/etc/service/syslog/run
#!/bin/sh
exec syslogd -n -l ${SYSLOG_LEVEL-6} -O -
