# reverse-proxy/root/etc/service/dnsmasq/run
#!/bin/sh
exec dnsmasq --keep-in-foreground --user=root --no-negcache
