# nginx/root/etc/service/nginx/finish
#!/bin/sh
if [ $1 -ne 0 -a $1 -ne 256 ]; then
    exec halt
fi
