# gitweb/root/etc/service/fcgiwrap/run
#!/bin/sh
mkdir -p /var/run/fcgiwrap
chown nginx:daemon /var/run/fcgiwrap
exec s6-setuidgid nginx fcgiwrap -s unix:/var/run/fcgiwrap/fcgiwrap.sock
