# pritunl/etc/service/pritunl/run
#!/bin/sh
pritunl set-mongodb mongodb://localhost:27017/pritunl
exec pritunl start
