# seafile/etc/service/seahub/run
#!/bin/sh
exec gunicorn seahub.wsgi:application -c /srv/seafile-server/runtime/seahub.conf -b 0.0.0.0:80 --enable-stdio-inheritance --capture-output
