# base/root/etc/service/crond/run
#!/bin/execlineb -P
/usr/sbin/crond -f -d1
