# base/root/etc/service/crond/run
#!/bin/sh
exec crond -f -d ${CRON_LOG_LEVEL-9}
